# webstatus.service
[Unit]
Description=WebStatus Network Monitoring System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/webstatus
Environment="PATH=/home/pi/webstatus/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/pi/webstatus/venv/bin/python /home/pi/webstatus/main.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webstatus

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
